# dist/linux/newrelic/init/nr-clearblockd-testnet.conf
#!upstart
description "clearblockd (testnet) (newrelic)"
limit nofile 32768 32768

# When to start the service
start on startup

# When to stop the service
stop on shutdown

# Automatically restart process if crashed
respawn
respawn limit 10 60 # 10 times in 60 seconds

script
USER=xcpd
USER_HOME=/home/xcp
NR_ADMIN_PATH=$USER_HOME/clearinghoused_build/env.clearblockd/bin/newrelic-admin
PYTHON_INTERP_PATH=$USER_HOME/clearinghoused_build/env.clearblockd/bin/python
PYTHON_SCRIPT_PATH=$USER_HOME/clearinghoused_build/dist/clearblockd/clearblockd.py
NEW_RELIC_CONFIG_FILE=/etc/newrelic/nr_clearblockd.ini
export NEW_RELIC_CONFIG_FILE

start-stop-daemon --start --pidfile $USER_HOME/.config/clearblockd-testnet/clearblockd.pid --chuid $USER --exec $NR_ADMIN_PATH -- run-program $PYTHON_INTERP_PATH $PYTHON_SCRIPT_PATH --data-dir=$USER_HOME/.config/clearblockd-testnet --testnet
end script
